# deploy/voice-agent.service
# Systemd service file for Voice Scheduling Agent
# Copy to: /etc/systemd/system/voice-agent.service
# Then run: sudo systemctl daemon-reload && sudo systemctl enable voice-agent

[Unit]
Description=Voice Scheduling Agent
Documentation=https://github.com/your-repo/voice-scheduling-agent
After=network.target
Wants=network-online.target

[Service]
Type=simple
User=ubuntu
Group=ubuntu
WorkingDirectory=/home/ubuntu/voice-scheduling-agent

# Environment file
EnvironmentFile=/home/ubuntu/voice-scheduling-agent/.env

# Python virtual environment
Environment="PATH=/home/ubuntu/voice-scheduling-agent/venv/bin:/usr/local/bin:/usr/bin:/bin"

# Run the application
ExecStart=/home/ubuntu/voice-scheduling-agent/venv/bin/python run.py

# Restart policy
Restart=always
RestartSec=10

# Logging
StandardOutput=append:/var/log/voice-agent/stdout.log
StandardError=append:/var/log/voice-agent/stderr.log

# Security hardening
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=/home/ubuntu/voice-scheduling-agent/data /home/ubuntu/voice-scheduling-agent/logs

# Resource limits
LimitNOFILE=65535
LimitNPROC=4096

[Install]
WantedBy=multi-user.target
